# scripts/timekpr-app.service
[Unit]
Description=timekpr-app web interface
After=network.target timekpr.service
Requires=timekpr.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/torgeir/prosjekter/timeKpr-app/timeKpr-app
Environment=PATH=/home/torgeir/prosjekter/timeKpr-app/timeKpr-app/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/torgeir/prosjekter/timeKpr-app/timeKpr-app/src
EnvironmentFile=-/home/torgeir/prosjekter/timeKpr-app/timeKpr-app/.env
ExecStart=/home/torgeir/prosjekter/timeKpr-app/timeKpr-app/.venv/bin/python -m timekpr_app.api.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
